# systemd/project-builder-webhook.service
[Unit]
Description=Project Builder Webhook Server
After=network.target

[Service]
Type=simple
User=projectbuilder
Group=projectbuilder
WorkingDirectory=/home/projectbuilder

# Main service command
ExecStart=/usr/bin/python3 /opt/project-builder/webhook-server.py

# Restart configuration
Restart=always
RestartSec=10
StartLimitBurst=5

# Environment variables
Environment=PATH=/usr/local/bin:/usr/bin:/bin
EnvironmentFile=-/etc/environment

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=project-builder-webhook

[Install]
WantedBy=multi-user.target